# alloy/rootfs/etc/s6-overlay/s6-rc.d/init-alloy/run
#!/usr/bin/with-contenv bash
set -e

CONFIG_DIR="/etc/alloy"
CONFIG_FILE="${CONFIG_DIR}/config.alloy"
OPTIONS_FILE="/data/options.json"

mkdir -p "${CONFIG_DIR}" /data/alloy

# Read options
LOKI_URL=$(jq -r '.loki_url' "${OPTIONS_FILE}")
# Trim surrounding whitespace — easy to paste into the HA options UI and
# otherwise crash-loops Alloy with a cryptic URL parse error.
LOKI_URL=$(printf '%s' "${LOKI_URL}" | sed -e 's/^[[:space:]]*//' -e 's/[[:space:]]*$//')
LOG_LEVEL=$(jq -r '.log_level // "info"' "${OPTIONS_FILE}")
JOURNAL_PRIORITY_AS_LEVEL=$(jq -r '.journal_priority_as_level | if . == null then true else . end' "${OPTIONS_FILE}")
PARSE_HA_LOG_LEVEL=$(jq -r '.parse_ha_log_level // false' "${OPTIONS_FILE}")
ADDITIONAL_CONFIG=$(jq -r '.additional_config // ""' "${OPTIONS_FILE}")

# Optional journal-priority -> level relabel rule. Docker's journald
# driver assigns a fixed priority per stream (stdout=info, stderr=err),
# so for container logs this label reflects the stream, not the real
# level. Disable to leave the level to content-based detection.
if [ "${JOURNAL_PRIORITY_AS_LEVEL}" = "true" ]; then
    LEVEL_RULE='  rule {
    source_labels = ["__journal_priority_keyword"]
    target_label  = "level"
  }'
else
    LEVEL_RULE=""
fi

# Optional pipeline stages extracting the Python application log level
# from HA Core / Supervisor lines into an ha_level label. Scoped via
# stage.match so other journal streams are unaffected.
HA_LEVEL_STAGES=''
if [ "${PARSE_HA_LOG_LEVEL}" = "true" ]; then
    HA_LEVEL_STAGES='
  stage.match {
    selector = "{container_name=~\"homeassistant|hassio_supervisor\"}"

    stage.regex {
      expression = "\\d{4}-\\d{2}-\\d{2}[ T]\\d{2}:\\d{2}:\\d{2}[.,]\\d+\\s+(?P<ha_level>DEBUG|INFO|WARNING|ERROR|CRITICAL)"
    }

    stage.labels {
      values = { ha_level = "" }
    }
  }
'
fi

# Detect journal path
JOURNAL_PATH="/var/log/journal"
if [ ! -d "${JOURNAL_PATH}" ] || [ -z "$(ls -A "${JOURNAL_PATH}" 2>/dev/null)" ]; then
    JOURNAL_PATH="/run/log/journal"
fi

echo "-----------------------------------------------------------"
echo " Grafana Alloy for Home Assistant"
echo " Loki URL: ${LOKI_URL}"
echo " Journal path: ${JOURNAL_PATH}"
echo " Log level: ${LOG_LEVEL}"
echo " Journal priority as level: ${JOURNAL_PRIORITY_AS_LEVEL}"
echo " Parse HA log level: ${PARSE_HA_LOG_LEVEL}"
echo "-----------------------------------------------------------"

# Generate Alloy config
cat > "${CONFIG_FILE}" <<ALLOYCONFIG
// Grafana Alloy config - auto-generated by HAOS add-on
// Do not edit directly - configure via add-on options

logging {
  level = "${LOG_LEVEL}"
}

// --- Read systemd journal ---
loki.source.journal "journal" {
  path         = "${JOURNAL_PATH}"
  forward_to   = [loki.process.journal.receiver]
  relabel_rules = loki.relabel.journal.rules
  labels       = {
    job = "systemd-journal",
  }
}

// --- Relabel journal fields ---
loki.relabel "journal" {
  forward_to = []

  rule {
    source_labels = ["__journal__systemd_unit"]
    target_label  = "unit"
  }
  rule {
    source_labels = ["__journal__hostname"]
    target_label  = "hostname"
  }
  rule {
    source_labels = ["__journal_syslog_identifier"]
    target_label  = "syslog_identifier"
  }
  rule {
    source_labels = ["__journal__transport"]
    target_label  = "transport"
  }
  rule {
    source_labels = ["__journal_container_name"]
    target_label  = "container_name"
  }
${LEVEL_RULE}
}

// --- Process journal logs ---
loki.process "journal" {
  stage.drop {
    expression = "^\\\\s*$"
  }
${HA_LEVEL_STAGES}
  forward_to = [loki.write.loki.receiver]
}

// --- Write to Loki ---
loki.write "loki" {
  endpoint {
    url = "${LOKI_URL}"
  }
}
ALLOYCONFIG

# Append additional config if provided
if [ -n "${ADDITIONAL_CONFIG}" ] && [ "${ADDITIONAL_CONFIG}" != "null" ]; then
    echo "" >> "${CONFIG_FILE}"
    echo "// --- Additional user config ---" >> "${CONFIG_FILE}"
    echo "${ADDITIONAL_CONFIG}" >> "${CONFIG_FILE}"
fi

echo "Alloy config generated at ${CONFIG_FILE}"
